# supervisord.conf
[supervisord]
nodaemon=false
logfile=/home/user/webapp/supervisord.log
pidfile=/home/user/webapp/supervisord.pid

[unix_http_server]
file=/home/user/webapp/supervisor.sock

[supervisorctl]
serverurl=unix:///home/user/webapp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:anomaly_agent_api]
command=python api_server.py
directory=/home/user/webapp
environment=OPENAI_API_KEY="YOUR_OPENAI_API_KEY_HERE"
autostart=true
autorestart=true
stdout_logfile=/home/user/webapp/api_server.log
stderr_logfile=/home/user/webapp/api_server_error.log
user=user
startsecs=10
startretries=3